# dev-indexer_1/systemd/spool-watcher.service
[Unit]
Description=Spool Directory Watcher for RAG Ingestion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/zenglow-indexer
WorkingDirectory=%h/ZenGlow/ZenGlow/dev-indexer_1
ExecStart=/usr/bin/env python3 scripts/spool_poll_watcher.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
